# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=xbindkeys
Exec=xbindkeys -f ~/.xbindkeysrc
Icon=keyboard
