# .config/autostart/xcompmgr.desktop
[Desktop Entry]
Type=Application
Name=xcompmgr
Exec=xcompmgr -c
Icon=computerjanitor
